# scripts/.config/systemd/user/sleepprompt.service
[Unit]
Description=System Sleep with prompt

[Service]
Type=oneshot
ExecStart=/home/noel/bin/sleepprompt
